# agentforge.service
# === PURE RUST FLYWHEEL DEFAULT CUTOVER (Phase 3 — PRODUCTION EXCELLENCE) @ 2026-05-31T10:42:02+03:00 ===
# Injected by bin/make_pure_rust_flywheel_default.sh
# agentforge-runner sole engine (flywheel-step/candidate/continuous). Non-breaking additive.
# Rollback: edit unit or set Environment=AGENTFORGE_FLYWHEEL_ENGINE=python + touch .disable_pure_rust_flywheel
[Unit]
Description=AgentForge Task Queue Server
After=network.target

[Service]
Type=simple
User=agx
WorkingDirectory=/home/eveselove/agentforge
# PHASE 4 / task-5af0e350: API now Rust gateway (agentforge-gateway binary). task_queue.py .bak only.
# Use dedicated agentforge-gateway.service for clean gw; this kept for ref/compat.
ExecStart=/home/eveselove/agentforge/gateway/target/release/agentforge-gateway
Restart=always
RestartSec=5
Environment=PATH=/home/eveselove/.cargo/bin:/home/eveselove/.grok/bin:/home/eveselove/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/eveselove
# task-5af0e350: PURE RUST FLYWHEEL HARD DEFAULT (soak timer) - applies to all flywheel paths even if this unit is gateway.
Environment=AGENTFORGE_PURE_RUST_FLYWHEEL=1
Environment=AGENTFORGE_FLYWHEEL_ENGINE=rust
Environment=FLYWHEEL_PROVENANCE=rust-agentforge-runner
Environment=AGENTFORGE_FLYWHEEL_PROVENANCE=rust-agentforge-runner

[Install]
WantedBy=multi-user.target

# Pure cutover + task-5af0e350: PURE RUST HARD DEFAULT for 14d soak. Rollback via disable_pure script.
# PHASE4_REMOVAL_PLAN AGGRESSIVE FINAL DEPRECATION SWEEP: flywheel mentions are for rollback+env only; unit updated to gw binary.
