# agentforge-watchdog.service
# === PURE RUST FLYWHEEL DEFAULT CUTOVER (Phase 3 — PRODUCTION EXCELLENCE) @ 2026-05-31T10:42:02+03:00 ===
# Injected by bin/make_pure_rust_flywheel_default.sh
# agentforge-runner sole engine (flywheel-step/candidate/continuous). Non-breaking additive.
# Rollback: edit unit or set Environment=AGENTFORGE_FLYWHEEL_ENGINE=python + touch .disable_pure_rust_flywheel
[Unit]
Description=AgentForge Watchdog Daemon
After=agentforge-gateway.service network.target agentforge-api.service

[Service]
Type=simple
# PHASE 4: Watchdog (Python monitor for health). Flywheel via runner.
Environment=AGENTFORGE_RUST_FLYWHEEL=1
Environment=AGENTFORGE_USE_RUST=1
Environment=AGENTFORGE_RUST_RUNNER=/home/eveselove/agentforge/rust/target/release/agentforge-runner
ExecStart=/usr/bin/python3 /home/eveselove/agentforge/watchdog.py
WorkingDirectory=/home/eveselove/agentforge
Restart=always
RestartSec=5
StandardOutput=append:/home/eveselove/agentforge/logs/watchdog.log
StandardError=append:/home/eveselove/agentforge/logs/watchdog.log

[Install]
WantedBy=default.target

Environment=AGENTFORGE_PURE_RUST_FLYWHEEL=1
Environment=AGENTFORGE_FLYWHEEL_ENGINE=rust
Environment=AGENTFORGE_FLYWHEEL_PROVENANCE=rust-agentforge-runner
# Pure cutover injected by make_pure_rust_flywheel_default.sh @ 2026-05-31T10:42:02+03:00 — rollback via FLYWHEEL_ENGINE=python or .disable_pure_rust_flywheel
